# roles/music/files/sound.service
[Unit]
Description=Sound

[Service]
ExecStart=/bin/bash /usr/bin/sound.sh

[Install]
WantedBy=multi-user.target